# vpn-guard.service
[Unit]
Description=VPN Guard Agent monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vpn-guard top --config /etc/vpn-guard/config.yml
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
